# conf/ihatemoney.conf
[program:__APP__]
command=__FINALPATH__/venv/bin/gunicorn -c /etc/__APP__/gunicorn.conf.py __APP__.wsgi:application
user=__APP__
autostart=true
autorestart=true
redirect_stderr=true
